# childcare360.service
[Unit]
Description=Childcare360 Childcare Platform v1.2.0
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/childcare-roster-app
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3003
EnvironmentFile=-/root/childcare-roster-app/.env

[Install]
WantedBy=multi-user.target
