# childcare360.service
[Unit]
Description=Childcare360 Application Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/childcare360-app
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=5
StandardOutput=append:/root/childcare360-app/childcare360.log
StandardError=append:/root/childcare360-app/childcare360.log
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
